# deploy/paper-bot-15m-boll.service
[Unit]
Description=Hyperliquid fade paper bot (15m, Bollinger z>=2.5 — Phase 3 trigger arm)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hyper
WorkingDirectory=/opt/hyperdata
ExecStart=/usr/bin/python3 /opt/hyperdata/paper_bot.py --interval 15m --trigger bollinger --datadir /opt/hyperdata/paper_15m_boll
Restart=always
RestartSec=15
StandardOutput=append:/opt/hyperdata/paper_15m_boll/service.out
StandardError=append:/opt/hyperdata/paper_15m_boll/service.err

[Install]
WantedBy=multi-user.target
