# deploy/paper-bot-5m.service
[Unit]
Description=Hyperliquid fade paper bot (5m)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hyper
WorkingDirectory=/opt/hyperdata
EnvironmentFile=-/etc/hyperdata/telegram.env
ExecStart=/usr/bin/python3 /opt/hyperdata/paper_bot.py --interval 5m --datadir /opt/hyperdata/paper_5m
Restart=always
RestartSec=15
StandardOutput=append:/opt/hyperdata/paper_5m/service.out
StandardError=append:/opt/hyperdata/paper_5m/service.err

[Install]
WantedBy=multi-user.target
